# conf/supervisor.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/flask/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx -g "daemon off;"
